# recipes-apps/camera-kiosk/files/camera-kiosk.service
[Unit]
Description=IMX219 Camera Kiosk Mode - Direct DisplayPort Output
After=nvargus-daemon.service systemd-udev-settle.service
Requires=nvargus-daemon.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/camera-kiosk.sh
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
