# recipes-apps/camera-kiosk/files/camera-kiosk.service
[Unit]
Description=IMX219 Camera Kiosk Mode - Fullscreen Display
After=weston-kiosk.service nvargus-daemon.service
Requires=weston-kiosk.service nvargus-daemon.service

[Service]
Type=simple
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="WAYLAND_DISPLAY=wayland-0"
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/camera-kiosk.sh
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
